# templates/etc/profile.d/hardening.sh.j2
# {{ ansible_managed }}

# Setting umask for all users
umask {{ hardening_login_umask }}

# Setting timeout if no input is given to shell
TMOUT={{ hardening_login_timeout }}
